# com.infosecfollow.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.infosecfollow.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/mlac/AI_Projects/infosecfollow/run_daily.sh</string>
  </array>
  <!-- 6 AM, 12 PM, 4 PM, 9 PM in the Mac's local (Eastern) time. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/mlac/AI_Projects/infosecfollow/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mlac/AI_Projects/infosecfollow/logs/launchd.err.log</string>
</dict>
</plist>
