# deploy/crontab
# supercronic schedule. Times are in the container's $TZ (set to America/New_York
# in docker-compose.yml), matching the original macOS LaunchAgent: 6am, noon, 4pm,
# and 9pm ET. Fields: minute hour day-of-month month day-of-week.
2 6  * * * /app/run-briefing.sh
2 12 * * * /app/run-briefing.sh
2 16 * * * /app/run-briefing.sh
2 21 * * * /app/run-briefing.sh
